# docker/siber.supervisord.conf
[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9012        ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=5           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true             ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9012 ; use an http:// url to specify an inet socket

[program:api-test-siber]
user=works
command=/home/works/program/api-test-siber --port=19080 --gwport=9080 --conf=/home/works/program/conf/online.conf -log=send
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
directory=/home/works/program/             ; directory to cwd to before exec (def no cwd)
autostart=true                ; start at supervisord start (default: true)
autorestart=true
startsecs=10                  ; number of secs prog must stay running (def. 1)
startretries=5                ; max # of serial start failures (default 3)
stdout_logfile=/home/works/program/logs/api-test-siber.sup.log       ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile=/home/works/program/logs/api-test-siber.sup.error        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)

[eventlistener:supervisor-monitor]
command=/bin/supervisor-monitor -c /etc/supervisor-event-listener.ini
; 指定配置文件路径
events=PROCESS_STATE_EXITED
